# bin/smf/manifest.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<service_bundle type="manifest" name="export">
  <service name="site/prox" type="service" version="1">

    <create_default_instance enabled="false"/>

    <single_instance/>

    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>

    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <dependency name="mongodb" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/pkgsrc/mongodb:default"/>
    </dependency>

    <method_context working_directory="/usr/local/prox">
      <method_credential user="admin" group="staff" privileges="basic,net_privaddr"  />
      <method_environment>
        <envvar name="PATH" value="/home/admin/local/bin:/usr/local/bin:/usr/bin:/usr/sbin:/bin"/>
        <envvar name="HOME" value="/home/admin"/>
      </method_environment>
    </method_context>

    <exec_method
      type="method"
      name="start"
      exec="/usr/local/prox/bin/prox start"
      timeout_seconds="5"/>

    <exec_method
      type="method"
      name="stop"
      exec="/usr/local/prox/bin/prox stop"
      timeout_seconds="60"/>

    <exec_method
      type="method"
      name="restart"
      exec="/usr/local/prox/bin/prox restart"
      timeout_seconds="60"/>

    <property_group name="startd" type="framework">
      <propval name="critical_failure_count" type="astring" value="3"/>
      <propval name="critical_failure_period" type="astring" value="10"/>
      <propval name="duration" type="astring" value="transient"/>
      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>

    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">Proxibase Service</loctext>
      </common_name>
    </template>

  </service>

</service_bundle>
